# ops/systemd/medtech-hr-erp-backup.timer
[Unit]
Description=Run MedTech HR ERP backups hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true
Unit=medtech-hr-erp-backup.service

[Install]
WantedBy=timers.target
